# conf.d/brav.conf
[program:brav]
command=/app/brav.sh
autorestart=false
